# systemd/user/cerberus-web-dashboard.service
[Unit]
Description=CERBERUS Web Command Center Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ph/cerberus-honeypot
ExecStart=/home/ph/cerberus-honeypot/venv/bin/python3 scripts/web_dashboard.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cerberus-web-dashboard

[Install]
WantedBy=default.target
